# nexus-pipeline/jenkins/supervisord.conf
[supervisord]
nodaemon=true

[program:jenkins]
priority=10
command=/bin/tini -- /usr/local/bin/jenkins.sh

[program:python]
priority=20
command=python /tmp/upload-pipelines.py
