# jenkins/dedicated-node-docker/supervisord.conf
[supervisord]
nodaemon=true

[program:docker]
command=dockerd-entrypoint.sh

[program:node]
command=python /home/node-attach.py
