# supervisord/supervisord.conf
[supervisord]
nodaemon=true

[program:openssh]
command=/usr/sbin/sshd -D

[program:pythonscript]
command=python /home/python_script.py