# cloudbees/operation_center/supervisord.conf
[supervisord]
nodaemon=true

[program:jenkins_oc]
command=/bin/tini -- /usr/local/bin/jenkins.sh

[program:activation_script]
command=python /home/configure_trial_license.py